# deploy/assetframe-daily.service
[Unit]
Description=AssetFrame daily scheduled run (score closed windows + generate all due assets)
Documentation=https://github.com/JoeWat2005/assetframe-scripts
After=network-online.target
Wants=network-online.target

[Service]
# oneshot: the timer fires this, it runs the daily production batch, then exits.
Type=oneshot
# --- EDIT these three paths to match the VM (same as the poller unit) ---
WorkingDirectory=/opt/assetframe-scripts
EnvironmentFile=/opt/assetframe-scripts/.env
ExecStart=/opt/assetframe-scripts/.venv/bin/python -m scripts.scheduler.scheduled_run
# ------------------------------------------------------------------------

Environment=PYTHONUNBUFFERED=1
# A full BATCH run (score + author/critique every due asset via the async Message-Batches API) can
# legitimately run 1.5-2h+, because batch latency is queue-dependent. 3h ceiling gives comfortable
# headroom so the run never gets SIGKILL'd mid-publish (which orphans the engine_runs row at
# 'running'). Keep ASSETFRAME_RUN_TIMEOUT (engine.json) <= this so the engine's own kill fires first
# and records a clean status. Lower both to ~1800 if you switch to the synchronous path.
TimeoutStartSec=10800

NoNewPrivileges=true
ProtectSystem=full

# No [Install] — this unit is driven by assetframe-daily.timer, not enabled directly.
