# deploy/assetframe-control.service
[Unit]
Description=AssetFrame box control server (HTTP API behind cloudflared + Cloudflare Access) — PROD
Documentation=https://github.com/JoeWat2005/assetframe-scripts
# Outbound only — binds 127.0.0.1; the public hostname engine.assetframe.co.uk reaches it ONLY
# through the cloudflared tunnel, gated by Cloudflare Access. Start after the tunnel.
After=network-online.target cloudflared-engine.service
Wants=network-online.target

[Service]
Type=simple
# --- EDIT these paths to match the VM (same checkout as the prod poller) ---
WorkingDirectory=/opt/assetframe-scripts
EnvironmentFile=/opt/assetframe-scripts/.env
ExecStart=/opt/assetframe-scripts/.venv/bin/python -m scripts.coordination.control_server
# --------------------------------------------------------------------------
# The prod .env MUST set:  ASSETFRAME_CONTROL_PORT=8787
#                          ASSETFRAME_CF_TEAM=wat044
#                          ASSETFRAME_CONTROL_AUD=<prod Access app AUD>
#               (REQUIRED)  ASSETFRAME_CONTROL_TOKEN=<random secret; the POST /control write factor —
#                           ALSO a Vercel env var so the dashboard sends it as Authorization: Bearer>
# Without CF_TEAM + CONTROL_AUD + CONTROL_TOKEN the server refuses to start (fail-closed), unless
# ASSETFRAME_CONTROL_INSECURE=1 is set (LOCAL DEV ONLY — never on the box).
# OPTIONAL: ASSETFRAME_POLLER_UNIT names the poller unit checked for `online` via `systemctl
#           is-active` (default assetframe-poller; the DEV box must set assetframe-poller-dev). This
#           replaces the heartbeat as the online signal, so the poller writes no Upstash heartbeat.

Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=15
Environment=PYTHONUNBUFFERED=1

# Runs as ROOT (no User=) on purpose: Restart engine / Deploy latest code call `systemctl restart
# assetframe-poller` (+ self-restart) so they work even when the poller is DEAD. If you run this as a
# non-root user instead, grant it `systemctl restart assetframe-poller assetframe-control` via sudoers
# or polkit, else those two commands fail (everything else still works).
# Light hardening: it needs outbound HTTPS (Neon + the team JWKS), to read+write the repo tree (deploy),
# and dbus to systemd (restart). ProtectSystem=full keeps /usr,/boot,/etc read-only; /opt stays writable.
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=false

[Install]
WantedBy=multi-user.target
